# ansible/ansible-practice/playbooks/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=node_exporter
ExecStart=/opt/node_exporter
Restart=always

[Install]
WantedBy=multi-user.target